# System_node.service
[Unit]
Description=System Node
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/System_node/System_node.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target